# launchd/com.regeorge.stickergotaro.xiaoman-nightly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.regeorge.stickergotaro.xiaoman-nightly</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/regeorge/Documents/codeStore/stickergotaro/.venv-xiaogpt-local/bin/python</string>
      <string>/Users/regeorge/Documents/codeStore/stickergotaro/scripts/run_nightly_pipeline.py</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>23</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/regeorge/Documents/codeStore/stickergotaro/data/logs/xiaoman-nightly.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/regeorge/Documents/codeStore/stickergotaro/data/logs/xiaoman-nightly.stderr.log</string>

    <key>RunAtLoad</key>
    <false/>
  </dict>
</plist>
